# templates/default/logstash.upstart.erb
# Logstash <%= node['logstash']['component'] %>
#
# Generated by Chef

description "logstash <%= node['logstash']['component'] %>"

start on runlevel [2345]
stop on runlevel [06]

# tell upstart we're creating a daemon
# upstart manages PID creation for you.
expect fork

script
  cd <%= node['logstash']['install_path'] %>
  exec su -c "java -jar logstash-<%= node['logstash']['version'] %>-monolithic.jar web --log <%= node['logstash']['log_path']/logstash-web.log &" - <%= node['logstash']['user'] %>
  emit logstash-web-server_running
end script
